# wifiproxy.service
[Unit]
Description=WiFi Proxy Web UI
After=network-online.target dnsmasq.service
Wants=network-online.target

[Service]
Type=simple
User=wifiproxy
WorkingDirectory=/opt/wifiproxy
Environment="BIND_HOST=192.168.50.1"
Environment="PORT=80"
Environment="SECRET_KEY=CHANGE_ME_BEFORE_DEPLOY"
ExecStart=/opt/wifiproxy/venv/bin/python run.py
Restart=on-failure
RestartSec=5

# Allow binding to port 80 as a non-root user.
# AmbientCapabilities and NoNewPrivileges=yes are mutually exclusive —
# the capability must be ambient (inherited), not elevated.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
